# provision/programs.ini
[program:udp2raw]
directory = /vagrant
command = udp2raw -c --conf-file config/udp2raw.conf -a -l 127.0.0.1:3336 --disable-color
redirect_stderr = true

[program:kcptun]
command = kcptun -l 127.0.0.1:3335 -r 127.0.0.1:3336 --nocomp --crypt none --mode fast2  --conn 32 --mtu 1200
redirect_stderr = true

[program:ss-local]
command = ss-local -s 127.0.0.1 -p 3335 -b 10.0.2.15 -l 1080 -m aes-128-gcm -k foobarx -v
redirect_stderr = true

[program:ss-redir]
command = ss-redir -s 127.0.0.1 -p 3335 -b 127.0.0.1 -l 20034 -m aes-128-gcm -k foobarx -v
redirect_stderr = true

[program:privoxy]
command = privoxy --no-daemon /vagrant/provision/privoxy.conf
